# telegram_bot/config/supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=/app/bot.py
process_name=bot
user=bot
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

